# supervisord.conf
[supervisord]
nodaemon=true

[program:opsdash-serverd]
command=/usr/sbin/opsdash-server -f
directory=/usr/sbin
